# football-betting-pipeline/com.football.betting.run_final.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.football.betting.run_final</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhiwuzou/app/football-betting/football-betting-pipeline/.venv/bin/python</string>
    <string>run_final.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/zhiwuzou/app/football-betting/football-betting-pipeline</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/zhiwuzou/app/football-betting-log/football-betting-run-final.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zhiwuzou/app/football-betting-log/football-betting-run-final.err</string>
</dict>
</plist>
